# trunk/emsi/asp-net-app/PACRAT/noninteractive/.crontab
#
# $Id$
#
#
# Monthly
#
25 22 01  * * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/PACRAT/noninteractive/statement_service_strike_team_roster.aspx
#
# Special
#
30 22 01 01 * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/PACRAT/noninteractive/statement_strike_team_member_status.aspx
30 22 01 02 * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/PACRAT/noninteractive/statement_strike_team_member_status.aspx?do_limit_to_uncredentialed=true
30 22 01 03 * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/PACRAT/noninteractive/statement_strike_team_member_status.aspx?do_limit_to_uncredentialed=true
30 22 01 04 * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/PACRAT/noninteractive/statement_strike_team_member_status.aspx
30 22 01 05 * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/PACRAT/noninteractive/statement_strike_team_member_status.aspx?do_limit_to_uncredentialed=true
30 22 01 06 * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/PACRAT/noninteractive/statement_strike_team_member_status.aspx?do_limit_to_uncredentialed=true
30 22 01 07 * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/PACRAT/noninteractive/statement_strike_team_member_status.aspx
30 22 01 08 * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/PACRAT/noninteractive/statement_strike_team_member_status.aspx?do_limit_to_uncredentialed=true
30 22 01 09 * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/PACRAT/noninteractive/statement_strike_team_member_status.aspx?do_limit_to_uncredentialed=true
30 22 01 10 * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/PACRAT/noninteractive/statement_strike_team_member_status.aspx
30 22 01 11 * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/PACRAT/noninteractive/statement_strike_team_member_status.aspx?do_limit_to_uncredentialed=true
30 22 01 12 * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/PACRAT/noninteractive/statement_strike_team_member_status.aspx?do_limit_to_uncredentialed=true
