# asp-net-app/WebEMSOF/noninteractive/.crontab
#
# $Id$
#
@daily wget -q https://sslws010.alentus.com/frompaper2web/WebEMSOF/noninteractive/milestone_sweep.aspx
