# trunk/emsi/asp-net-app/WebEMSOF/noninteractive/.crontab
#
# $Id$
#
00 00 * * * wget -q https://sslws010.alentus.com/frompaper2web/WebEMSOF/noninteractive/milestone_sweep.aspx