# trunk/emsi/asp-net-app/WebEMSOF/noninteractive/.crontab
#
# $Id$
#
mm hh DD MM d (mm=minute, hh=hour, DD=DayOfMonth, MM=MonthNum, d=DayOfWeekWhereZeroIsSunday)
00 00  *  * * wget --output-document=/dev/null https://sslws010.alentus.com/frompaper2web/WebEMSOF/noninteractive/milestone_sweep.aspx
00 01  1  6 * wget --output-document=/dev/null https://sslws010.alentus.com/frompaper2web/WebEMSOF/noninteractive/report_equipment_procurement_overview.aspx
