# trunk/emsi/asp-net-app/ConEdLink/noninteractive/.crontab
#
# $Id$
#
#mmhh DD MM d (mm=minute, hh=hour, DD=DayOfMonth, MM=MonthNum, d=DayOfWeekWhereZeroIsSunday)
#
# Daily
#
50 03  *  * * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/ConEdLink/noninteractive/daily_tasks.aspx
#
# Weekly
#
30 19  *  * 2 wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/ConEdLink/noninteractive/import_latest_practitioners_from_emsrs.aspx
#
# Yearly
#
00 20 01 01 * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/ConEdLink/noninteractive/first_of_calendar_month_tasks.aspx
